# includes/etc/supervisor/conf.d/supervisor.conf
[program:apache2]
command=/start-apache.sh
numprocs=1
autostart=true
autorestart=true

[program:mysqld]
command=/start-mysql.sh
numprocs=1
autostart=true
autorestart=true